# systemd/user/bt-audio-autoswitch.service
[Unit]
Description=Auto-switch audio sink when Pixel Buds connect/disconnect
After=pipewire.service wireplumber.service
Wants=pipewire.service wireplumber.service

[Service]
Type=simple
Environment=PATH=/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/bin/bash %h/.local/bin/bt-audio-autoswitch.sh
Restart=always
RestartSec=2

[Install]
WantedBy=default.target

